# httpd/httpd.ini
[supervisord]
nodaemon=true

[program:httpd]
command=/usr/bin/httpd -D FOREGROUND
